# remote/.local/share/systemd/user/remote-upload@.service
[Unit]
Description=Upload media to cloud storage

[Service]
Type=oneshot
ExecStart=%h/.local/bin/media-upload

[Install]
WantedBy=timers.target