# remote/.local/share/systemd/user/remote-sync@.service
[Unit]
Description=Remote Cloud Syncronisation

[Service]
Type=oneshot
ExecStart=%h/.local/bin/remote sync "%i"

[Install]
WantedBy=timers.target